# systemd/system/tcsd.service.d/anti-evil-maid-seal.conf
# start early enough for anti-evil-maid-seal.service

[Unit]
DefaultDependencies=false
# for tcsd_changer_identify:
Requires=local-fs.target
After=local-fs.target
