# trousers-changer/systemd/system/tcsd.service.d/tcsd_changer_identify.conf
[Service]
ExecStartPre=/sbin/tcsd_changer_identify
